# crontab.cron
0 * * * * bash  /Users/neils/SemEnv/ifconfig2Docs.sh
