# crontab.cron
0 * * * * bash  /Users/neilstevens/SemEnv/ifconfig2Docs.sh
